# srcs/supervisord.conf
[supervusord]
nodaemon=true

[program:sshd]
command = sh -C "/usr/bin/sshd -D && kill -s SIGTERM $(cat supervisord.pid)"

[program:nginx]
command = sh -C "nginx -g 'daemon off;' && kill -s SIGTERM $(cat supervisord.pid)"